# test/portTest.service
[Unit]
Description=My Custom Startup Program
After=network.target
 
[Service]
Type=simple
ExecStart=/path/to/your/ELFfile
 
[Install]
WantedBy=multi-user.target